# deploy/linux/systemd/finex-auth.service
[Unit]
Description=FinEx Auth Service
After=network.target

[Service]
Type=simple
User=finex
Group=finex
WorkingDirectory=/opt/finex/backend
ExecStart=/bin/bash -lc 'set -a; source /opt/finex/backend/.env.production; set +a; exec /usr/bin/java ${FINEX_AUTH_JAVA_OPTS:--Xms256m -Xmx512m} -jar /opt/finex/backend/auth-service.jar'
Restart=always
RestartSec=5
SuccessExitStatus=143
StandardOutput=append:/opt/finex/logs/auth-service.log
StandardError=append:/opt/finex/logs/auth-service.err.log

[Install]
WantedBy=multi-user.target
